# config/.config/systemd/user/radio.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/icecast2 -c icecast.xml

[Install]
WantedBy=default.target
